# mdnsd.service
[Unit]
Description=mdns daemon
[Service]
ExecStart=/bin/mdnsd

[Install]
WantedBy=multi-user.target

