# mdnsd.service
[Unit]
Description=mdns daemon
StartLimitIntervalSec=0
[Service]
ExecStart=/bin/picomdnsd


[Install]
WantedBy=multi-user.target

